# load_test/systemd/gifgloo-loadtest-api.service
[Unit]
Description=Gifgloo loadtest API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gifgloo-backend-loadtest
EnvironmentFile=/home/ubuntu/gifgloo-backend-loadtest/.env.loadtest
Environment=PYTHONUNBUFFERED=1
Environment=PROMETHEUS_MULTIPROC_DIR=/tmp/gifgloo-loadtest-prometheus
Environment=LOADTEST_API_WORKERS=2
ExecStartPre=/usr/bin/install -d -o ubuntu -g ubuntu -m 0755 /tmp/gifgloo-loadtest-prometheus
ExecStartPre=/usr/bin/find /tmp/gifgloo-loadtest-prometheus -mindepth 1 -maxdepth 1 -type f -name *.db -delete
ExecStart=/home/ubuntu/gifgloo-backend-loadtest/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8001 --workers ${LOADTEST_API_WORKERS}
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=multi-user.target
